# resources/CPack_scripts/postinst
#!/bin/sh
systemctl daemon-reload
sleep 1
systemctl enable my.service
systemctl start my.service